# crontab/cron_certbot-renew
# Attempt to renew SSL certificates daily
# Best practice is to set a different sleep timer for different servers
0 4 * * * sleep 120 && /root/.local/bin/bin-server/certbot-renew >> /root/certbot-renew.log
